# tools/systemd-user/organize-downloads.path
[Unit]
Description=Watch Downloads for new files

[Path]
PathChanged=/home/guayaba/Downloads
Unit=organize-downloads.service

[Install]
WantedBy=default.target
